# config/supervisord.conf
[supervisord]
logfile=/var/log/supervisord.log                ; supervisord log file
logfile_maxbytes=50MB                           ; maximum size of logfile before rotation
logfile_backups=10                              ; number of backed up logfiles
loglevel=error                                  ; info, debug, warn, trace
pidfile=/var/run/supervisord.pid                ; pidfile location
nodaemon=true                                   ; run supervisord NOT as a daemon
minfds=1024                                     ; number of startup file descriptors
minprocs=200                                    ; number of process descriptors
user=root                                       ; default user
childlogdir=/var/log/                           ; where child log files will live

[program:app-gunicorn]
command = /usr/local/bin/gunicorn --chdir /app/ bctf.wsgi -b 0.0.0.0:8000
redirect_stderr=true


[program:nginx]
command = /usr/sbin/nginx
redirect_stderr=true

[program:crond]
command = /usr/sbin/crond -b
redirect_stderr=true
